# src/supervisord.conf
# suporvisord.conf
[program:uwsgi]
environment=PATH='/opt/conda/envs/ecsdemo/bin'
command=/opt/conda/envs/ecsdemo/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini --die-on-term --need-app